# tarot-bot.service
[Unit]
Description=Tarot Bot Service
After=network.target

[Service]
Type=simple
User=tarotbot
Group=tarotbot
WorkingDirectory=/opt/tarot_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
